# deploy/systemd/openapi-mcp.service
[Unit]
Description=OpenAPI MCP Builder (Trimble Agentic AI Platform)
Documentation=https://github.com/Charley-Forey-AI/openapi-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mcp
Group=mcp
WorkingDirectory=/mcp/openapi-mcp
EnvironmentFile=/mcp/openapi-mcp/.env
ExecStart=/mcp/openapi-mcp/.venv/bin/openapi-mcp-builder
Restart=on-failure
RestartSec=3
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openapi-mcp

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ReadWritePaths=/mcp/openapi-mcp
AmbientCapabilities=
CapabilityBoundingSet=
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
